# build/linux/ondevice-daemon.service
[Unit]
Description = ondevice daemon

[Service]
User = ondevice
Environment = CONFFILE=/var/lib/ondevice/ondevice.conf
Environment = PIDFILE=/var/run/ondevice/ondevice.pid
Environment = SOCKFILE=/var/run/ondevice/ondevice.sock
EnvironmentFile = /etc/default/ondevice-daemon
ExecStart = /usr/bin/ondevice daemon --conf="${CONFFILE}" --pidfile="${PIDFILE}" --sock="${SOCKFILE}"

RestartSec = 5
Restart = on-failure
